# gitspace_uwsgi.conf
description "gitspace uwsgi"

start on runlevel [2345]
stop on runlevel [06]

exec uwsgi /home/jeremydw/gitspace/server.ini
